# root/etc/services.d/monit/run
#!/usr/bin/with-contenv bash
exec \
    s6-setuidgid abc /usr/bin/monit -I -B -c /config/monitrc